# wifi-watchdog.timer
[Unit]
Description=Periodic wlan0 supplicant health check

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
